# cydia/iForward/Library/LaunchDaemons/com.iforward.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" 
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">

<dict>
    <key>Label</key>
    <string>com.iforward</string>

    <key>UserName</key>
    <string>mobile</string>
    
    <key>GroupName</key>
    <string>admin</string>
    
    <key>Program</key>
    <string>/usr/bin/iForward</string>

   <key>RunAtLoad</key>
   <true/>

    <key>StartInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/var/mobile/Library/Logs/iForward.log</string>

    <key>StandardErrorPath</key>
    <string>/var/mobile/Library/Logs/iForward.log</string>

</dict>
</plist>
